# altimeafactory/supervisord.conf
[supervisord]
nodaemon=true

[program:nonversioned]
command=/deploy.sh
numprocs=1
autostart=true
autorestart=unexpected
exitcodes=0